# deploy/systemd/ffball.service
[Unit]
Description=Fantasy Football Draft Board
After=network.target docker.service
Requires=docker.service
Wants=postgresql.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/app
Environment=PATH=/nix/var/nix/profiles/default/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/nix/var/nix/profiles/default/bin/docker info
ExecStart=/nix/var/nix/profiles/default/bin/docker compose up -d
ExecStop=/nix/var/nix/profiles/default/bin/docker compose down
TimeoutStartSec=300
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target 